# dev/vagrant/puppet/templates/etc/init.d/frontend.erb
#!/bin/bash
#
# <%= @role %>  Frontend webapp
#
# chkconfig:    2345 80 20
# description:  Frontend webapp
# processname:  <%= @role %>
# pidfile:      /executable-jar-apps/<%= @role %>/<%= @role %>.pid
#
# source function library
. /etc/rc.d/init.d/functions

USER=jvmuser
APP_HOME=/executable-jar-apps/<%= @role %>
PIDFILE=${APP_HOME}/RUNNING_PID
RETVAL=0

start() {
  if [ -e "${APP_HOME}/<%= @role %>.jar" ]
  then
    echo -n $"Starting <%= @role %>: "
    pushd ${APP_HOME} > /dev/null
    daemon --user ${USER} --pidfile ${PIDFILE} \
      java -jar <%= @role %>.jar 2>&1 >> logs/application.log &
    RETVAL=$?

    # sleep for a second to ensure the process has started before we terminate
    # the controlling TTY (magenta deployment scripts need this otherwise the
    # container doesn't start)
    sleep 1

    echo
    popd > /dev/null
  else
    echo "Not started: ${APP_HOME}/<%= @role %>.jar not present."
  fi
}

stop() {
  echo -n $"Stopping <%= @role %>: "
  killproc -p $PIDFILE <%= @role %>
  RETVAL=$?
  echo
}
   
restart() {
  stop
  start
}

case "$1" in
  start)
    start
    ;;

  stop)
    stop
    ;;

  restart|force-reload|reload)
    restart
    ;;

  status)
    status -p ${PIDFILE} <%= @role %>
    RETVAL=$?
    ;;

  *)
    echo $"Usage: $0 {start|stop|status|restart|reload|force-reload}"
    exit 1

esac
  
exit $RETVAL